# tilemill.conf
[program:tilemill]
user=root
command=/run_tilemill.sh
autorestart=true
stopsignal=INT
